# plain/setup/rc.local
#!/bin/bash

openvt -s -w /media/Vault/plain/bin/onboot-init

exit 0
